# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D -e
autorestart=true
startretries=3

[program:cloud-init-start]
command=/usr/sbin/cloud-init-start.sh
